# packaging/razortrade-signals.timer
# razortrade Donchian signal timer — fires every 4 hours
#
# Cadence rationale: the strategy uses 4h Donchian candles, so the
# natural evaluation frequency is once per candle close. Candles close
# on 4h UTC boundaries (00/04/08/12/16/20). We schedule on those exact
# boundaries plus a randomised jitter so we don't hammer Kraken at
# :00:00 sharp alongside every other bot.
#
# Persistent=true ensures that if the host was offline at a scheduled
# time (reboot, network outage, maintenance window) the timer fires
# the backlog tick when the host returns — important for not missing a
# breakout evaluation.

[Unit]
Description=Run razortrade Donchian signal generator every 4 hours
Documentation=file:/usr/share/doc/razortrade/README.md
Requires=razortrade-signals.service

[Timer]
# Every 4h: 00:00, 04:00, 08:00, 12:00, 16:00, 20:00 UTC — exactly at
# Kraken's 4h candle close boundaries. The explicit UTC suffix is
# REQUIRED; without it, systemd interprets OnCalendar in the host's
# local timezone (CEST on our prod box), causing the timer to fire 2h
# late relative to candle close and shifting by another hour twice a
# year on DST transitions. This matters because the Donchian signal
# script pulls 4h candles and evaluates breakout conditions — firing
# late means we act on signals with added latency, firing on DST
# transitions means we fire at :00:30 CEST on the shift weekend.
#
# Do not drop the "UTC" suffix. Previous version pre-Drop-17 shipped
# without it and the timer ran on CEST; symptoms were subtle (no
# obvious errors) which is why it survived paper validation.
OnCalendar=*-*-* 00/4:00:00 UTC

# Randomised delay up to 2 minutes. Spreads load and survives the
# "every bot fires at :00:00" thundering herd at Kraken.
RandomizedDelaySec=120

# If we miss a scheduled fire because the host was down, run it on
# next boot instead of silently skipping. Keeps the signal history
# gap-free in the expected case.
Persistent=true

Unit=razortrade-signals.service

[Install]
WantedBy=timers.target
